# config/etc/supervisor/conf.d/nodejs.conf
[program:nodejs]
directory=/app
command=/usr/bin/npm start
stopsignal=TERM
redirect_stderr=true
autorestart=true
