# 12/etc/services.d/app/run
#!/usr/bin/execlineb -P

with-contenv
cd /app
/usr/bin/npm start
